# .config/agent-slack/com.kush.agent-slack-fleet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.kush.agent-slack-fleet</string>
  <key>ProgramArguments</key>
  <array><string>__HOME__/.local/bin/agent-slack</string><string>fleet</string></array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>EnvironmentVariables</key>
  <dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:__HOME__/.local/bin</string></dict>
  <key>StandardOutPath</key><string>/tmp/agent-slack-fleet.log</string>
  <key>StandardErrorPath</key><string>/tmp/agent-slack-fleet.log</string>
</dict>
</plist>
